# service-monitor.service
[Unit]
Description=Service timeout monitor

[Service]
Type=simple
ExecStart=/usr/bin/service-monitor.py

[Install]
WantedBy=basic.target
